# debian/gnome-clips.postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        if command -v systemctl >/dev/null 2>&1; then
            systemctl --global enable gnome-clips-daemon.service || true
        fi
        ;;
esac

#DEBHELPER#
exit 0
